# src/srv/tangd.socket.in
[Unit]
Description=Tang daemon socket
ConditionDirectoryNotEmpty=@TANG_DB@

[Socket]
ListenStream=@TANG_PORT@
Accept=true

[Install]
WantedBy=multi-user.target
